# pipelines/pipeline-1/dashboard-nextjs/cloudflared.service
[Unit]
Description=Cloudflare Tunnel for My Life in Data Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lliverani
ExecStart=/usr/bin/cloudflared tunnel --config /home/lliverani/.cloudflared/config.yml run mylife-dashboard
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
