# pipelines/pipeline-1/dashboard-nextjs/mylife-dashboard.service
[Unit]
Description=My Life in Data - Next.js Dashboard
After=network.target

[Service]
Type=simple
User=lliverani
WorkingDirectory=/home/lliverani/projects/mylife-in-data/pipelines/pipeline-1/dashboard-nextjs
Environment=NODE_ENV=production
# First build the app, then start it
ExecStartPre=/usr/bin/npm run build
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
